# linux/debian_family/system_service_modifications.yml
- id: 1021
  dataTypes:
    - linux
  name: System Service Modifications
  impact:
    confidentiality: 2
    integrity: 4
    availability: 3
  category: System Integrity
  technique: T1543.001 - Create or Modify System Process
  adversary: origin
  references:
    - https://attack.mitre.org/techniques/T1543/001/
    - https://www.debian.org/doc/manuals/debian-reference/ch03.en.html
  description: Detects modifications to system services including systemd units, init scripts, or service configurations that could indicate persistence attempts or system compromise
  where: |
    (safe(log.process_name, "") == "systemctl" && 
     (safe(log.command_line, "") contains "enable" || 
      safe(log.command_line, "") contains "disable" || 
      safe(log.command_line, "") contains "mask" || 
      safe(log.command_line, "") contains "unmask")) ||
    (safe(log.file_path, "") contains "/etc/systemd/system/" && 
     safe(log.event_type, "") == "file_modify") ||
    (safe(log.file_path, "") contains "/etc/init.d/" && 
     safe(log.event_type, "") == "file_modify") ||
    (safe(log.message, "") contains "Created symlink" && 
     safe(log.message, "") contains "/etc/systemd/system/")
  afterEvents:
    - indexPattern: v11-log-*
      with:
        - field: origin.user.keyword
          operator: filter_term
          value: '{{origin.user}}'
      within: now-15m
      count: 3
  deduplicateBy:
    - origin.user
    - log.service_name